# embedded_sensor_service/systemd/assignment_sensor.service
[Unit]
Description=Mock sensor logger (assignment)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/assignment_sensor --interval 5 --logfile /tmp/assignment_sensor.log
Restart=on-failure
RestartSec=2
# Opcional: correr como usuario no-root (ajustar permisos si se usa)
# User=someuser

[Install]
WantedBy=multi-user.target
